# launchd/com.plei99.piano-tracker.sync.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.plei99.piano-tracker.sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>__TRACKER_BIN__</string>
		<string>sync</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>StartInterval</key>
	<integer>3600</integer>

	<key>StandardOutPath</key>
	<string>__STDOUT_LOG__</string>

	<key>StandardErrorPath</key>
	<string>__STDERR_LOG__</string>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
